# setup-ai-process/design-library-reindex.service
[Unit]
Description=Full re-index of design library
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=oneshot
User=rpi
Group=rpi
WorkingDirectory=/home/rpi/ai-engine/design-library-indexer
Environment="CHROMA_TELEMETRY=FALSE"
Environment="ANONYMIZED_TELEMETRY=FALSE"
Environment="LIBRARY_ROOT=/mnt/design-library"
Environment="CHROMA_DIR=/mnt/design-library/chroma_data"
EnvironmentFile=-/home/rpi/ai-engine/.env
ExecStart=/home/rpi/ai-engine/venv/bin/python run_indexer.py index --full \
    --library-root ${LIBRARY_ROOT} \
    --chroma-dir ${CHROMA_DIR}
StandardOutput=append:/home/rpi/ai-engine/logs/reindex.log
StandardError=append:/home/rpi/ai-engine/logs/reindex-error.log

# Generous timeout — full index can take hours on Pi
TimeoutStartSec=28800

MemoryMax=1G
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=-/mnt/design-library/.index -/mnt/design-library/chroma_data -/home/rpi/ai-engine/design-library/.index -/home/rpi/ai-engine/chroma_data /home/rpi/ai-engine/logs
